# launchd/com.roffe.smartpdf.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.roffe.smartpdf.watch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/roffe/src/smartpdf/smartpdf-watch.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>WatchPaths</key>
  <array>
    <string>/Users/roffe/Dropbox/From_BrotherDevice</string>
  </array>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>EnvironmentVariables</key>
  <dict>
    <key>SMARTPDF_SOURCE_DIR</key>
    <string>/Users/roffe/Dropbox/From_BrotherDevice</string>
    <key>SMARTPDF_DEST_DIR</key>
    <string>/Users/roffe/Tresorit/Scanned Documents</string>
    <key>SMARTPDF_TRASH_DIR</key>
    <string>/Users/roffe/.Thrash</string>
    <key>SMARTPDF_BIN</key>
    <string>/Users/roffe/src/smartpdf/smartpdf.sh</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/roffe/Library/Logs/smartpdf-watch.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/roffe/Library/Logs/smartpdf-watch.log</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
